# roles/macos/player/templates/com.orangead.player.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{ player.service_name }}</string>
    
    <key>Program</key>
    <string>{{ player.video_source }}/scripts/start_player.sh</string>
    
    <key>WorkingDirectory</key>
    <string>{{ player.video_source }}</string>
    
    <key>StandardOutPath</key>
    <string>{{ player.video_source }}/logs/player.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{ player.video_source }}/logs/player.log</string>
    
    <!-- Log consolidation: Using a single log file for all output -->
    
    <key>RunAtLoad</key>
    <{{ player.service_run_at_load | default(true) | lower }}/>
    
    <key>KeepAlive</key>
    <{{ player.service_keep_alive | default(true) | lower }}/>
    
    {% if player.restart_on_failure | default(true) %}
    <key>ThrottleInterval</key>
    <integer>{{ player.restart_delay | default(10) }}</integer>
    {% endif %}
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{ ansible_user_dir }}</string>
        <key>USER</key>
        <string>{{ ansible_user }}</string>
        <key>DISPLAY</key>
        <string>:0</string>
        <key>SDL_VIDEODRIVER</key>
        <string>cocoa</string>
    </dict>
    
    <!-- Ensure the service can access the display -->
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>
    
    <!-- Process type for GUI access -->
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <!-- Set nice level for video playback performance -->
    <key>Nice</key>
    <integer>-5</integer>
    
    <!-- Disable low priority I/O -->
    <key>LowPriorityIO</key>
    <false/>
    
    <!-- Additional LaunchAgent settings for reliability -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
    
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>
    
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>512</integer>
    </dict>
</dict>
</plist>