# roles/common/device_api/templates/com.orangead.deviceapi.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.orangead.deviceapi</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>{{ device_api_python }}</string>
        <string>{{ device_api_path }}/main.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>{{ device_api_path }}</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>{{ device_api_path }}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>{{ device_api_path }}/logs/deviceapi.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{ device_api_path }}/logs/deviceapi_error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LaunchOnlyOnce</key>
    <false/>
    
    <key>UserName</key>
    <string>{{ ansible_user }}</string>
</dict>
</plist>