# roles/ubuntu/server_optimization/templates/ethernet-pm-disable.service.j2
# {{ ansible_managed }}
# Ethernet Power Management Disable Service
# Ensures ethernet interfaces maintain full performance and stability

[Unit]
Description=Disable Ethernet Power Management for OrangeAd Infrastructure
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
{% if ubuntu_ethernet_optimization.power_management.disable_wol %}
# Disable Wake-on-LAN for all ethernet interfaces
{% for interface in ethernet_interfaces %}
ExecStart=/usr/sbin/ethtool -s {{ interface }} wol d
{% endfor %}
{% endif %}

{% if ubuntu_ethernet_optimization.power_management.disable_interface_pm %}
# Additional power management settings can be added here
# Example: ExecStart=/usr/sbin/ethtool -s eth0 msglvl 0
{% endif %}

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

# Applied on: {{ ansible_date_time.iso8601 }}
# By: oaAnsible ubuntu/server_optimization role