# roles/macos/network/tailscale/templates/com.tailscale.tailscaled.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tailscale.tailscaled</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ go_bin_path_for_tailscale }}/tailscaled</string>
        <string>--state=/var/lib/tailscale/tailscaled.state</string>
        <string>--socket=/var/run/tailscaled.socket</string>
        <string>--port=41641</string>
        <!-- <string>--verbose=1</string> --> <!-- Uncomment for more verbose tailscaled logging if needed -->
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <!-- Keep the service running. If it exits for any reason other than a clean exit (status 0),
             launchd will restart it. Tailscaled is designed to handle its own retry logic
             for network connectivity to the control plane. -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Additionally, if the socket file disappears (indicating a crash or unexpected stop),
             launchd will also restart it. -->
        <key>PathState</key>
        <dict>
            <key>/var/run/tailscaled.socket</key>
            <false/>
        </dict>
    </dict>
    <!-- Helps ensure all related processes are cleaned up if the main daemon is stopped/restarted -->
    <key>AbandonProcessGroup</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/tailscaled.out.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/tailscaled.err.log</string>
    <key>WorkingDirectory</key>
    <string>/var/lib/tailscale</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>UserName</key>
    <string>root</string>
</dict>
</plist>