# roles/macos/alpr_service/templates/com.orangead.alpr.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{ macos_alpr_service_plist }}</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
        # Comprehensive container startup script with cleanup
        set -e
        
        # Remove any existing container with the same name
        /usr/local/bin/docker rm -f {{ macos_alpr_container_name }} 2>/dev/null || true
        
        # Wait a moment for cleanup
        sleep 2
        
        # Start the container with proper configuration
        exec /usr/local/bin/docker run \
            --name {{ macos_alpr_container_name }} \
            --restart=unless-stopped \
            --rm \
{% if ansible_machine == 'arm64' %}
            --platform linux/arm64 \
{% endif %}
            -p {{ macos_alpr_host_port }}:{{ macos_alpr_container_port }} \
            -v {{ macos_alpr_license_volume_name }}:/license \
            -e TOKEN={{ vault_alpr_token }} \
            -e LICENSE_KEY={{ vault_alpr_license_key }} \
            {{ macos_alpr_image_name }}
        </string>
    </array>
    
    <key>RunAtLoad</key>
    <{{ macos_alpr_run_at_load | lower }}/>
    
    <key>KeepAlive</key>
    <{{ macos_alpr_keep_alive | lower }}/>
    
    <key>StandardOutPath</key>
    <string>{{ ansible_user_dir }}/orangead/alpr/logs/alpr.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{ ansible_user_dir }}/orangead/alpr/logs/alpr.err.log</string>
    
    <key>WorkingDirectory</key>
    <string>{{ ansible_user_dir }}/orangead/alpr</string>
    
    <key>UserName</key>
    <string>{{ ansible_user }}</string>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{ homebrew_prefix }}/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{ ansible_user_dir }}</string>
    </dict>
</dict>
</plist>